# cloudflared-pi.service
[Unit]
Description=Cloudflare Tunnel - pi
After=network.target

[Service]
Type=simple
User=youruser
ExecStart=/home/youruser/.local/bin/cloudflared tunnel --config /home/youruser/.cloudflared/config-pi.yml run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
