# cloudflared-pi.user.service
[Unit]
Description=Cloudflare Tunnel - pi (User Service)
After=network.target

[Service]
Type=simple
ExecStart=/home/youruser/.local/bin/cloudflared tunnel --config /home/youruser/.cloudflared/config-pi.yml run
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
