# pika.user.service
[Unit]
Description=Pika - AI Coding Companion (User Service)
After=network.target cloudflared-pi.service

[Service]
Type=simple
WorkingDirectory=/home/leo/code/pi-agent-manager
ExecStart=/home/leo/code/pi-agent-manager/target/release/pika
Restart=always
RestartSec=5
Environment=PORT=7847
Environment=PATH=/home/leo/.nvm/versions/node/v22.18.0/bin:/usr/local/bin:/usr/bin:/bin

# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
