# deploy/pi-agent-manager.service
[Unit]
Description=Pi Agent Manager
After=network.target cloudflared-pi-api.service

[Service]
Type=simple
User=leo
WorkingDirectory=/home/leo/code/pi-agent-manager
ExecStart=/home/leo/code/pi-agent-manager/target/release/pi-agent-manager
Restart=always
RestartSec=5
Environment=PORT=7847

# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
